# var/lib/dpkg/info/zerotier-one.conffiles
/etc/init.d/zerotier-one
/etc/init/zerotier-one.conf
